# src/service/test-service.service.in
[Unit]
Description=SomeIP test application

[Service]
ExecStart=@prefix@/bin/common-api-someip-test-service
LimitCORE=infinity

[Install]
